# gfx-hat-stats.service
[Unit]
Description=GFX HAT System Stats Display
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi
ExecStart=/usr/bin/python3 /home/dietpi/gfx_hat_stats.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
